# ngsquid/supervisord.conf
[supervisord]  
nodaemon=true  
  
[program:squid]  
command=/sbin/squidentrypoint.sh
  
[program:nginx]  
command=/usr/local/openresty/bin/openresty -g "daemon off;"